# DataCollector/crontab
0 10 * * * ~/Dropbox/ALMA/WeatherStation/DataCollector/
